# udev-rules/81-garmin-usb-sync.rules
ACTION=="add", SUBSYSTEM=="usb", SYSFS{idVendor}=="091e", SYSFS{idProduct}=="27af", RUN+="/home/jrbenito/bin/garmin-sync.sh"
